# etc/services.d/named/run
#!/usr/bin/execlineb -P

named -u abc -g
